# supervisord.conf
[supervisord]
nodaemon=true

[program:myrun]
priority=10
command=/usr/bin/myrun
autorestart=false
stdout_logfile=/var/log/myrun.log
